# backend/scripts/com.nexus.resolve-events.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.resolve-events</string>

    <!-- Configure NEXUS_BASE_DIR before loading: launchctl setenv NEXUS_BASE_DIR /path/to/Nexus-setup -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>node "$NEXUS_BASE_DIR/scripts/resolve-raw-events.js"</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NEXUS_HOST</key>
        <string>${NEXUS_DB_HOST}</string>
        <key>NEXUS_PORT</key>
        <string>${NEXUS_DB_PORT:-5432}</string>
        <key>NEXUS_DB</key>
        <string>${NEXUS_DB_NAME:-nexus}</string>
        <key>NEXUS_USER</key>
        <string>${NEXUS_DB_USER:-nexus}</string>
        <key>NEXUS_PASSWORD</key>
        <string>${NEXUS_DB_PASSWORD}</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/nexus-resolve-events.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/nexus-resolve-events.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/tmp</string>
</dict>
</plist>
